# systemd/nvidia-gddr6-fan-control.service
[Unit]
Description=NVIDIA GPU fan control driven by GDDR6/GDDR6X VRAM junction temp
# Start after Tailscale so --listen-tailscale can resolve the IP at boot.
After=network-online.target tailscaled.service
Wants=network-online.target tailscaled.service

[Service]
Type=simple
ExecStart=/usr/local/bin/nvidia-gddr6-fan-control --listen-tailscale
Restart=on-failure
RestartSec=5
# stdout/stderr are captured by journald automatically.

[Install]
WantedBy=multi-user.target
